# hubbot.service
[Unit]
Description=Hubbot

[Service]
ExecStart=/usr/local/bin/hubbot +
User=hubbot

[Install]
WantedBy=multi-user.target
