# hubbot.service
[Unit]
Description=Hubbot

[Service]
ExecStart=/usr/local/bin/hubbot +
